# com.claude-telegram-remote.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-telegram-remote.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>INSTALL_DIR/start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOGS_DIR/bot.log</string>

    <key>StandardErrorPath</key>
    <string>LOGS_DIR/error.log</string>

    <key>WorkingDirectory</key>
    <string>INSTALL_DIR</string>
</dict>
</plist>
